# deploy/systemd/arduino-trader-trading.service
[Unit]
Description=Arduino Trader - Trading Service
After=network.target

[Service]
Type=simple
User=arduino
WorkingDirectory=/home/arduino/arduino-trader
Environment="PYTHONPATH=/home/arduino/arduino-trader"
Environment="DEVICE_CONFIG_PATH=/home/arduino/arduino-trader/app/config/device.yaml"
ExecStart=/home/arduino/arduino-trader/venv/bin/python services/trading/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
