# systemd/sentinel-deploy.service
[Unit]
Description=Sentinel Auto-Deploy

[Service]
Type=oneshot
User=arduino
ExecStart=/home/arduino/sentinel/scripts/auto-deploy.sh
